# version1/setup1/hproducer-vm5/OutputMeter-config/crontab
PATH=/usr/local/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/custom/bin
* * * * * /bin/bash -c /custom/bin/execute-hproducer.sh
