# version1/setup1/eproducer-vm2/SmartMeter-config/crontab
SHELL=/bin/bash
* * * * * /custom/bin/execute.sh
* * * * * touch /abc.txt
